# build/etc/systemd/system/audio-control.service
[Unit]
Description=repeater audio control

[Service]
ExecStart=/opt/repeater-audio-control/scraper
WorkingDirectory=/opt/repeater-audio-control
Restart=always

[Install]
WantedBy=multi-user.target
