# build/etc/systemd/system/repeater-audio-control-backend.service
[Unit]
Description=repeater audio control

[Service]
ExecStart=/opt/repeater-audio-control/backend -config ./config.json
WorkingDirectory=/opt/repeater-audio-control
Restart=always

[Install]
WantedBy=multi-user.target
